# usbGadget.service
[Unit]
Description=USB Gadget init service

[Service]
Type=exec
ExecStart=/boot/usbGadget.sh

[Install]
WantedBy=multi-user.target

